# doc/cron/crontab
* * * * * djanog-admin update_index -a 1
0 * * * * djanog-admin refresh_sitemap
0 * * * * djanog-admin linuxos_cron
50 2 * * * djanog-admin rebuild_index --noinput
